# db/cassandra/crontab
* * * * * cqlsh -u aksuser -p password123 --execute=copy azuredata.azurecost from '/mnt/blob/output/test.csv' with header = true